# config/systemd/user/tnc@.target
[Unit]
# Version 1.0.0
Description=TNC on port %I
Wants=direwolf@%i.service kissport@%i.service kissattach@%i.service


