# config/systemd/user/tnc@.target
[Unit]
# Version 1.0.1
Description=Direwolf and AX25 TNC %I
Wants=direwolf@%i.service kissport@%i.service kissattach@%i.service


